# production/crontab
* * * * * www-data bash /srv/sites/production.khanovaskola.cz/production/worker.sh
