# production/supervisor.conf
[program:khanovaskola-beta-daemon]
command=php /srv/sites/beta.khanovaskola.cz/www/index.php worker
autostart=true
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog
